# home/dot_config/containers/systemd/dev-box.container
[Unit]
Description=My General Purpose Development Container

[Container]
Image=dev-box.build
ContainerName=my-dev-box

RunInit=true

# Works with host more smoothly
Network=host
SecurityLabelDisable=true
UserNS=keep-id
Volume=%h:%h
Tmpfs=/tmp

AddDevice=/dev/kfd
AddDevice=/dev/dri
Ulimit=memlock=-1:-1
SeccompProfile=unconfined
PodmanArgs=--ipc host

# Persist the container's ssh host keys
Volume=dev-box-ssh-keys:/etc/ssh

[Install]
WantedBy=default.target
