# home/dot_config/containers/systemd/dev-box.container
[Unit]
Description=My General Purpose Development Container

[Container]
Image=dev-box.build
ContainerName=my-dev-box

Network=host
SecurityLabelDisable=true
UserNS=keep-id
Volume=%h:%h

[Install]
WantedBy=default.target
